# crontab
MAILTO=jamespstrachan@gmail.com
59 13 * * * /usr/local/bin/python3.4 /usr/src/app/quickbook.py
